# scripts/atlas-memory-sync/atlas-memory-sync.service
[Unit]
Description=Atlas Memory File Sync
After=network.target atlas-memory.service
Wants=atlas-memory.service

[Service]
Type=simple
WorkingDirectory=/root/.openclaw/atlas-memory-sync
ExecStart=/usr/bin/python3 /root/.openclaw/atlas-memory-sync/sync.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
